# txt-to-md.service
[Unit]
Description=Serviço de conversão automática de TXT para MD
After=network.target

[Service]
# Verifique se o caminho abaixo está correto para o seu usuário
ExecStart=/bin/bash %h/automacao-txt-md/converter.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
